# etc/systemd/system/start-dhclient.service
[Unit]
Description=start-dhclient.sh script

[Service]
ExecStart=/startup/scripts/start-dhclient.sh

[Install]
WantedBy=multi-user.target
